# configs/check_leagues.conf
[program:check_leagues]
command=bash /home/varline/bin/check_leagues.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/check_leagues.err.log
stdout_logfile=/var/log/check_leagues.out.log