# configs/event_checker.conf
[program:event_checker]
command=bash /home/varline/bin/event_checker.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/event_checker.err.log
stdout_logfile=/var/log/event_checker.out.log