# configs/spider_live.conf
[program:spider_live]
command=bash /home/varline/bin/spider_live.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/spider_live.err.log
stdout_logfile=/var/log/spider_live.out.log